# configuration/boards/ldo-orbitool-o2s/98-ldo-orbitool-o2s.rules
# Orbitool O2s
SUBSYSTEMS=="usb", ATTRS{idProduct}=="614e",  ATTRS{idVendor}=="1d50", ATTRS{serial}=="ldo-orbitool-o2s", ACTION=="add", SYMLINK+="ldo-orbitool-o2s", RUN+="/home/pi/printer_data/config/Ratical/scripts/klipper-mcu-added.sh"
