# files/ratical-ondemand.service
[Unit]
Description=Ratical on-demand wake (configurator + crowsnest via nginx auth_request)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/ratical-ondemand/ratical-ondemand-wake.py
Restart=on-failure
RestartSec=2
Environment=RATICAL_ONDEMAND_IDLE_SEC=900
Environment=RATICAL_ONDEMAND_START_TIMEOUT=120
# Restrict to loopback wake API only
Environment=RATICAL_ONDEMAND_HOST=127.0.0.1
Environment=RATICAL_ONDEMAND_PORT=3199

[Install]
WantedBy=multi-user.target
